# docker/parser-supervisor.conf
[supervisord]
nodaemon=true

[program:parser]
directory=/app
command=/usr/bin/poetry run scrapy crawl mainpage-spider
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs=75
